# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=php-fpm

[program:crond]
command=/usr/sbin/crond
autostart=true
autorestart=true
